# daemon/claude-relay.service
[Unit]
Description=Claude Code Telegram Relay
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__BUN_PATH__ run src/bot.ts
Restart=always
RestartSec=10

# Environment
Environment=PATH=__BUN_DIR__:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=__PROJECT_DIR__/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-relay

[Install]
WantedBy=multi-user.target

# SETUP: Run `bash setup/install-systemd.sh` to install automatically.
#
# Manual commands:
#   sudo systemctl daemon-reload
#   sudo systemctl enable claude-relay
#   sudo systemctl start claude-relay
#   sudo systemctl status claude-relay
#   journalctl -u claude-relay -f
